# scripts/com.pierre.cours-uness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Identifiant unique du service -->
    <key>Label</key>
    <string>com.pierre.cours-uness</string>

    <!-- Script à exécuter -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/pierrecastet/Desktop/Cours_UNESS/scripts/git_auto_push.sh</string>
    </array>

    <!-- Déclencheurs : tout changement dans Fiches/ ou avancement_fiches.html -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/pierrecastet/Desktop/Cours_UNESS/Fiches</string>
        <string>/Users/pierrecastet/Desktop/Cours_UNESS/avancement_fiches.html</string>
    </array>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/pierrecastet/Desktop/Cours_UNESS/logs/git_push.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/pierrecastet/Desktop/Cours_UNESS/logs/git_push.err</string>

    <!-- Redémarrer si crash -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
